# init/el/7/captagent.service
[Unit]
Description=SIP capture agent server daemon
After=syslog.target network.target auditd.service

[Service]
EnvironmentFile=/etc/sysconfig/captagent
ExecStart=/usr/local/captagent/bin/captagent -f $CFG_FILE -d
KillMode=process
Restart=Always
RestartSec=42s
Type=forking

[Install]
WantedBy=multi-user.target
